# container-files/etc/supervisor.d/hhvm.conf
[program:hhvm]
command = hhvm --mode server --config /etc/hhvm/php.ini
autorestart = true
stdout_logfile = /data/logs/hhvm.log
stderr_logfile = /data/logs/hhvm.log
